# dotfiles/.config/systemd/user/tmux-resurrect-snapshot.service
[Unit]
Description=Snapshot newest tmux-resurrect save into the dotfiles repo

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/tmux-resurrect-snapshot.sh
